# distribution/systemd/minsky-supervisor.target
# Minsky supervisor target — groups all minsky child services so they start /
# stop together. Following systemd's documented "target as supervisor group"
# convention; not a process supervisor in the Erlang/OTP sense (systemd itself
# is the supervisor primitive — see ARCHITECTURE.md § "Process supervision
# tree" and vision.md § "Pattern conformance index" row 4 for the partial-
# conformance declaration).
#
# Install: copy to ~/.config/systemd/user/ after envsubst on placeholders.
# Placeholders this template uses (substituted by setup.sh):
#   ${MINSKY_HOME}   — absolute path to the minsky repo checkout
#
# Activation:
#   systemctl --user daemon-reload
#   systemctl --user enable --now minsky-supervisor.target

[Unit]
Description=Minsky supervisor target — autonomous Claude Code 24/7 loop
# Stability features (default-on per vision.md § 6 "Stay alive"): all
# services in this Wants= list are part of Minsky's stability contract —
# their absence makes the daemon less robust. New stability features
# MUST be added here in the same PR that ships their unit file. Opt-out
# is per-feature env var documented in the service file's header.
Wants=minsky-tick-loop.service
Wants=minsky-budget-guard.service
Wants=minsky-watchdog.service
After=network-online.target

[Install]
WantedBy=default.target
