# distribution/systemd/minsky-tick-loop.service
# Minsky tick-loop ŌĆö the per-iteration scheduler that picks the next task,
# claims it, runs OMC, and ships the result. Restart-on-failure with
# exponential backoff so a systematic failure doesn't hammer the API.
#
# Pattern: OTP supervision behaviour, "transient" restart strategy
# (Armstrong, *Programming Erlang*, 2007). systemd's `Restart=on-failure` +
# `RestartSec` + `RestartSteps` is the closest-fitting POSIX equivalent.
# Conformance: partial ŌĆö see ARCHITECTURE.md ┬¦ "Process supervision tree"
# and vision.md ┬¦ "Pattern conformance index" row 4.
#
# Placeholders (substituted by setup.sh via envsubst):
#   ${MINSKY_HOME}   ŌĆö absolute path to the minsky repo checkout

[Unit]
Description=Minsky tick-loop ŌĆö picks the next task and ships it
PartOf=minsky-supervisor.target
After=network-online.target minsky-budget-guard.service
Requires=minsky-budget-guard.service

[Service]
Type=simple
WorkingDirectory=${MINSKY_HOME}
Environment=MINSKY_HOME=${MINSKY_HOME}
# Safe-rollout window closed 2026-05-04 ŌĆö Minsky-on-itself validated
# end-to-end (launchd plist paths fixed #143; macOS portability fixes
# #144; node-PATH-resolution #145; stdout flush #146). Production
# default is `ProcessSpawnStrategy` (real `claude --print` per
# iteration). To re-enter dry-run temporarily (e.g., for a chaos
# test or a fresh-install smoke), uncomment the Environment= line
# below, `pnpm dogfood` to reload, then re-comment it.
# Environment=MINSKY_TICK_DRY_RUN=1
# Daily-changelog umbrella opt-in (acceptance criterion (3) of
# `daily-changelog-for-humans`). Enables BOTH `runDaemon` seams the CLI
# constructs under this single env var: the changelog-mode `claude --print`
# author (per-day, gated by `## YYYY-MM-DD` H2 in CHANGELOG.md) AND the
# per-day metric-snapshot capture (`pnpm changelog:snapshot --date <date>`,
# gated by `.minsky/metric-snapshots/<date>.json` presence). Both gates are
# idempotent ŌĆö exactly one fire per UTC date ŌĆö so leaving this on across
# many iterations is safe by design (rule #2 ŌĆö file presence IS the
# "this happened" record). Drop the flag to revert to substrate-inert.
Environment=MINSKY_CHANGELOG_ENABLE=1
# Post-task CTO-audit umbrella opt-in (sub-step (d) of `post-task-cto-audit`).
# Constructs the `CtoAuditSeam` the daemon dispatches into after every
# successfully-completed iteration that shipped a real change (commit on the
# branch OR PR opened ŌĆö noop iterations are gated out by `shouldRunCtoAudit`).
# The audit spawns a separate `claude --print` invocation with a CTO-mode
# prompt header (`CTO_PROMPT_HEADER`) that asks for the single highest-leverage
# next task; output is committed to TASKS.md on a dedicated branch
# (`audit/<UTC-date>-<original-task-id>`) and a PR is opened with label
# `minsky:cto-audit` for operator review (sub-step (e)). The cap-1-per-task
# contract (sub-step (f)) is enforced by a file-backed lock at
# `<MINSKY_HOME>/.minsky/cto-audit-lock/<task-id>` that survives daemon restart.
# Per-iteration override: set `MINSKY_CTO_AUDIT=off` to skip a single fire
# without unwinding the seam (sub-step (c)). Operator runbook:
# `docs/post-task-cto-audit.md`. Pivot trigger (rule #9): >5 audit-PRs/day
# means tighten the brief; <1/week sustained 4 weeks means retire.
Environment=MINSKY_CTO_AUDIT_ENABLE=1
ExecStart=/bin/sh ${MINSKY_HOME}/distribution/systemd/run-tick-loop.sh

# Erlang-style "transient" restart: only on non-zero exit (let-it-crash
# discipline, rule #6). Backoff caps the API hammering on systematic failures.
Restart=on-failure
RestartSec=5
StartLimitIntervalSec=300
StartLimitBurst=10

# StandardOutput/StandardError ŌåÆ journal (queryable via journalctl --user)
StandardOutput=journal
StandardError=journal
SyslogIdentifier=minsky-tick-loop

[Install]
WantedBy=minsky-supervisor.target
